# personal_systemd_units/backup.timer
[Unit]
Description=run backup.sh weekly 

[Timer]
OnCalendar=Mon *-*-* 17:55:00
Persistent=true

[Install]
WantedBy=timers.target
