# personal_systemd_units/crontab
50  17  *   *  1  termite --hold -e /home/anjan/backup.sh
